# templates/echo.service
[Unit]
Description=echoservice

[Service]
Type=forking
ExecStart=/usr/bin/service-script

[Install]
WantedBy=multi-user.target
